# bootstrap/ansible/roles/common/files/backup-data.timer
[Unit]
Description=Run backup-data every 10 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
